# rfs-debian/debian/var/lib/dpkg/info/tftpd-hpa.conffiles
/etc/init.d/tftpd-hpa
/etc/init/tftpd-hpa.conf
