# rfs-debian/debian/var/lib/dpkg/info/tftpd-hpa.list
/.
/etc
/etc/init.d
/etc/init.d/tftpd-hpa
/etc/init
/etc/init/tftpd-hpa.conf
/usr
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/in.tftpd.8.gz
/usr/share/doc
/usr/share/doc/tftpd-hpa
/usr/share/doc/tftpd-hpa/README.security
/usr/share/doc/tftpd-hpa/copyright
/usr/share/doc/tftpd-hpa/changelog.gz
/usr/share/doc/tftpd-hpa/changelog.Debian.gz
/usr/share/doc/tftpd-hpa/README.Debian
/usr/share/doc/tftpd-hpa/README
/usr/share/doc/tftpd-hpa/examples
/usr/share/doc/tftpd-hpa/examples/sample.rules
/usr/sbin
/usr/sbin/in.tftpd
/usr/share/man/man8/tftpd.8.gz
